# scripts/crontab/prod
#
# !!AUTO-GENERATED!!  Edit scripts/crontab/make-crons.py instead.
#

MAILTO=cron-sumo@mozilla.com

HOME = /tmp

# Every minute!
* * * * * cd /data/www/support.mozilla.com/kitsune; /usr/bin/python26 manage.py cron collect_tweets
* * * * * cd /data/www/support.mozilla.com/kitsune; /usr/bin/python26 manage.py cron get_queue_status

# Every hour.
42 * * * * cd /data/www/support.mozilla.com/kitsune; /usr/bin/python26 manage.py cleanup

# Every 2 hours.
1 */2 * * * cd /data/www/support.mozilla.com/kitsune; /usr/bin/python26 manage.py cron calculate_related_documents

# Every 6 hours.
0 */6 * * * cd /data/www/support.mozilla.com/kitsune; /usr/bin/python26 manage.py update_product_details -q > /dev/null
10 */6 * * * cd /data/www/support.mozilla.com/kitsune; /usr/bin/python26 manage.py cron rebuild_kb
30 */6 * * * cd /data/www/support.mozilla.com/kitsune; /usr/bin/python26 manage.py cron get_customercare_stats
40 */6 * * * cd /data/www/support.mozilla.com/kitsune; /usr/bin/python26 manage.py cron purge_tweets
50 */6 * * * cd /data/www/support.mozilla.com/kitsune; /usr/bin/python26 manage.py cron cache_top_contributors

# Once per day.
0 16 * * * cd /data/www/support.mozilla.com/kitsune; /usr/bin/python26 manage.py cron reload_wiki_traffic_stats
40 1 * * * cd /data/www/support.mozilla.com/kitsune; /usr/bin/python26 manage.py cron update_weekly_votes

# Twice per week.
#05 01 * * 1,4 cd /data/www/support.mozilla.com/kitsune; /usr/bin/python26 manage.py cron update_weekly_votes

# Once per week.
21 03 * * 3 cd /data/www/support.mozilla.com/kitsune; /usr/bin/python26 manage.py purge_hashes

MAILTO=root
